# dl/tensorflow/tensorboard/tensorboard.conf
[program:tensorboard]
command=tensorboard --port 6001 --logdir=/output
